# hyperspace.service
[Unit]
Description=Hyperspace Server
Documentation=https://github.com/hyperspace-org/hyperspace
After=network.target

[Service]
ExecStart=hyperspace
RestartSec=3

[Install]
WantedBy=multi-user.target
